# scripts/com.volregime.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchDaemon: weekly refresh of the Fixed-Income Vol-Regime tracker on the Mac mini.
  Reuses the EXISTING `vcdigest` service user (already has SSH, Tailscale, uv) — no new
  user needed. System domain → runs headless (SSH only, no GUI login). Distinct label and
  schedule from com.vcdigest.weekly so the two jobs don't clash.

  Runs scripts/weekly_update.sh: re-pull the live feeds (MOVE, FRED yields, NY Fed dealers),
  rebuild every data/*.json, then `git commit && git push`. GitHub Pages redeploys itself.

  Time is LOCAL to the mac (Europe/Podgorica). Here: Monday 08:00 (an hour before the
  vc-digest job at 09:00). PATH is set explicitly because the daemon's default PATH is
  stripped and wouldn't find git/uv.

  Git push needs non-interactive auth for `vcdigest` — add its SSH key as a repo deploy key
  (write access) first; see docs/deploy-macmini.md §2.

  Install (needs sudo; you have owner rights):
    sudo cp /Users/vcdigest/projects/vol-regime-tracker/scripts/com.volregime.weekly.plist /Library/LaunchDaemons/
    sudo chown root:wheel /Library/LaunchDaemons/com.volregime.weekly.plist
    sudo chmod 644        /Library/LaunchDaemons/com.volregime.weekly.plist
    sudo launchctl bootstrap system /Library/LaunchDaemons/com.volregime.weekly.plist
  Remove:  sudo launchctl bootout system /Library/LaunchDaemons/com.volregime.weekly.plist
  Run now: sudo launchctl kickstart -k system/com.volregime.weekly
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.volregime.weekly</string>

  <key>UserName</key>
  <string>vcdigest</string>
  <key>GroupName</key>
  <string>staff</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/vcdigest/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/vcdigest/projects/vol-regime-tracker/scripts/weekly_update.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/vcdigest/projects/vol-regime-tracker</string>

  <key>StandardOutPath</key>
  <string>/Users/vcdigest/projects/vol-regime-tracker/logs/weekly.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/vcdigest/projects/vol-regime-tracker/logs/weekly.err.log</string>

  <!-- Once a week: Monday 08:00 local mac time (Weekday 1 = Mon). -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>1</integer>
    <key>Hour</key><integer>8</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>
  <key>KeepAlive</key>
  <false/>
</dict>
</plist>
